# deploy/brosil.service
[Unit]
Description=Brosil tracker
After=network.target

[Service]
WorkingDirectory=/opt/brosil
ExecStart=/opt/brosil/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8090
Restart=always
User=root

[Install]
WantedBy=multi-user.target
